# systemd/issue-generator.service.in
[Unit]
Description=Generate issue file for login session
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=@SBINDIR@/issue-generator

[Install]
WantedBy=default.target
